# JDIS/home/sys999/cron/STID.cron
00 08 * * 1-5 /JDIS/home/sys999/cron/STID.script
