# JDIS/PROG/CMN/cleanup.cron
0 2 * * * cp /dev/null /etc/wtmp; cp /dev/null /usr/adm/sulog;
